# etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh
::ctrlaltdel:/bin/umount -a -r
ttyGS0::respawn:/sbin/getty -n -l /bin/sh ttyGS0 9600 vt100
